# input/rootfs/etc/s6/chrony/run
#!/bin/sh
exec /usr/sbin/chronyd -f /etc/chrony/chrony.conf -m -P 99
