# input/rootfs/etc/s6/gpsd/run
#!/bin/sh
exec /usr/sbin/gpsd -N -n -s 115200 /dev/ttyAMA0 /dev/pps0
