# input/rootfs/etc/s6/gpsd-listener/run
#!/bin/sh
exec s6-setuidgid nobody /bin/gpsd-listener.py
