# input/rootfs/etc/s6/gpsd-listener/run
#!/bin/sh
exec s6-setuidgid nobody /usr/bin/gpsd-listener.py
