# input/rootfs/etc/init.d/chrony-stats
#!/sbin/openrc-run
# shellcheck shell=ash

command=/bin/chrony-stats.py
pidfile="/var/run/chrony-stats.pid"
export PROMETHEUS_METRICS_FILE=/tmp/node-exporter/chrony.prom

depend()
{
        need net
        after firewall
        use dns
}

start() {
        ebegin "Starting chrony-stats.py"
    
        start-stop-daemon -m --background --start --quiet \
                --exec "${command}" \
                --pidfile "${pidfile}"

        eend $? "Failed to start chrony-stats.py"
}

stop() {
        ebegin "Stopping chrony-stats.py"
        start-stop-daemon --stop --quiet \
                --pidfile "${pidfile}"
        eend $? "Failed to stop chrony-stats.py"
}
