# input/rootfs/etc/init.d/data_prepare
#!/sbin/openrc-run
# shellcheck shell=ash

depend()
{
    need localmount
}

start()
{
    ebegin Preparing persistent data
    /sbin/ab_resizedata

    # make sure /data is mounted
    mount -a

    mkdir -p /data/chrony && chown -R chrony:chrony /data/chrony
    mkdir -p /tmp/node-exporter && chmod 1777 /tmp/node-exporter

    mkdir -p /data/etc /data/home /data/etc/ssh
    if [ ! -f /data/etc/shadow ]
    then
        cp -p /etc/shadow.alpine-builder /data/etc/shadow
    fi

    mkdir -p /data/root
    chmod 700 /data/root
    if [ ! -f /data/root/.zshrc ]
    then
        cp -rf /etc/skel/. /data/root/
    fi

    mount -a

    eend 0
}
