# input/rootfs/etc/init.d/gpsd-listener
#!/sbin/openrc-run
# shellcheck shell=ash

command=/bin/gpsd-listener.py
pidfile="/var/run/gpsd-listener.pid"

depend()
{
        need net
        after firewall
        use dns
}

start() {
        ebegin "Starting gpsd-listener.py"
    
        start-stop-daemon -m --background --start --quiet \
                --exec "${command}" \
                --pidfile "${pidfile}"

        eend $? "Failed to start gpsd-listener.py"
}

stop() {
        ebegin "Stopping gpsd-listener.py"
        start-stop-daemon --stop --quiet \
                --pidfile "${pidfile}"
        eend $? "Failed to stop gpsd-listener.py"
}
