# input/rootfs/etc/init.d/phc2sys
#!/sbin/openrc-run
# shellcheck shell=ash

command=/sbin/phc2sys
pidfile="/var/run/phc2sys.pid"
cfgfile=/etc/ptp4l.conf
iface=eth0.2
export PROMETHEUS_METRICS_FILE=/tmp/node-exporter/phc2sys.prom

depend()
{
        need net
        after firewall
        use dns
}

start() {
        ebegin "Starting phc2sys"
    
        start-stop-daemon -m --background --start --quiet \
                --exec /bin/ptp4l-sync-monitor.py \
                --pidfile "${pidfile}" \
                -- "${command}" -f "${cfgfile}" -s "${iface}" -c CLOCK_REALTIME -E ntpshm -M 1 -w -l 7 -m -q

        chrt -a -p 55 "$(cat "${pidfile}")"

        eend $? "Failed to start phc2sys"
}

stop() {
        ebegin "Stopping phc2sys"
        start-stop-daemon --stop --quiet \
                --pidfile "${pidfile}"
        eend $? "Failed to stop phc2sys"
}
