# input/rootfs/etc/init.d/phc2sys
#!/sbin/openrc-run
# shellcheck shell=ash

command=/sbin/phc2sys
pidfile="/var/run/phc2sys.pid"
cfgfile=/etc/ptp4l.conf

depend()
{
        need net
        after firewall
        use dns
}

start() {
        ebegin "Starting phc2sys"
    
        start-stop-daemon -m --background --start --quiet \
                --exec "${command}" \
                --pidfile "${pidfile}" \
                -- -f "${cfgfile}" -a -r -E linreg

        chrt -a -p 60 "$(cat "${pidfile}")"

        eend $? "Failed to start phc2sys"
}

stop() {
        ebegin "Stopping phc2sys"
        start-stop-daemon --stop --quiet \
                --pidfile "${pidfile}"
        eend $? "Failed to stop phc2sys"
}
