# input/rootfs/etc/init.d/ptp4l
#!/sbin/openrc-run
# shellcheck shell=ash

command=/usr/sbin/ptp4l
pidfile="/var/run/ptp4l.pid"
cfgfile=/etc/ptp4l.conf
ifaces="-i eth0 -i eth0.2 -i eth0.3 -i eth0.4 -i eth0.5 -i eth0.6"

depend()
{
        need net
        after firewall
        use dns
}

start() {
        ebegin "Starting ptp4l"
    
        start-stop-daemon -m --background --start --quiet \
                --exec "${command}" \
                --pidfile "${pidfile}" \
                -- $ifaces -f "${cfgfile}"

        chrt -a -p 50 "$(cat "${pidfile}")"

        eend $? "Failed to start ptp4l"
}

stop() {
        ebegin "Stopping ptp4l"
        start-stop-daemon --stop --quiet \
                --pidfile "${pidfile}"
        eend $? "Failed to stop ptp4l"
}
