# input/rootfs/etc/init.d/ts2phc
#!/sbin/openrc-run
# shellcheck shell=ash

command=/sbin/ts2phc
pidfile="/var/run/ts2phc.pid"
cfgfile=/etc/ptp4l.conf
iface=eth0
export PROMETHEUS_METRICS_FILE=/tmp/node-exporter/ts2phc.prom

depend()
{
        need net
        after firewall
        use dns
}

start() {
        ebegin "Starting ts2phc"

        mkdir -p /tmp/node-exporter

        start-stop-daemon -m --background --start --quiet \
                --exec /bin/ptp4l-sync-monitor.py \
                --pidfile "${pidfile}" \
                -- "${command}" -c "${iface}" -s nmea -f "${cfgfile}" -m -q -l 7

        chrt -a -p 60 "$(cat "${pidfile}")"

        eend $? "Failed to start ts2phc"
}

stop() {
        ebegin "Stopping ts2phc"
        start-stop-daemon --stop --quiet \
                --pidfile "${pidfile}"
        eend $? "Failed to stop ts2phc"
}
