# input/rootfs/etc/s6/chrony-stats/run
#!/bin/sh

export PROMETHEUS_METRICS_FILE=/tmp/node-exporter/chrony.prom

exec s6-setuidgid nobody /usr/bin/chrony-stats.py
