# input/rootfs/etc/s6/update-tai-offset/run
#!/bin/sh
exec /usr/bin/update-tai-offset.py
