# input/rootfs/etc/s6/update-tai-offset/run
#!/bin/sh
exec /bin/update-tai-offset.py
