# input/rootfs/etc/s6/node-exporter/run
#!/bin/sh
exec s6-setuidgid prometheus /usr/bin/node_exporter --collector.time --collector.timex --collector.ntp --collector.textfile.directory=/tmp/node-exporter
