# input/rootfs/etc/s6/keepalived/run
#!/bin/sh
exec /usr/sbin/keepalived --dont-fork --use-file=/etc/keepalived/keepalived.conf
