# input/rootfs/etc/s6/snmpd/run
#!/bin/sh
exec /usr/sbin/snmpd -f -c /data/snmpd.conf -LSwd -Lf /dev/null
