# input/rootfs/etc/s6/ts2phc/run
#!/bin/sh

export PROMETHEUS_METRICS_FILE=/tmp/node-exporter/ts2phc.prom

exec /bin/ptp4l-sync-monitor.py /usr/sbin/ts2phc -c eth0 -s nmea -f /etc/ptp4l.conf -m -q -l 7
